# scripts/polymarket-tick-compressor.service
[Unit]
Description=Polymarket L2 Tick Compressor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/home/botuser/polymarket-bot
Environment=PATH=/home/botuser/polymarket-bot/.venv/bin:/usr/bin:/bin
ExecStartPre=/bin/mkdir -p /home/botuser/polymarket-bot/data/l2_book_live /home/botuser/polymarket-bot/data/l2_book_live/_state
ExecStart=/home/botuser/polymarket-bot/.venv/bin/python /home/botuser/polymarket-bot/scripts/live_tick_compressor.py --output-dir /home/botuser/polymarket-bot/data/l2_book_live --market-limit 200 --flush-rows 10000 --flush-seconds 300 --rotation hourly --compression zstd --universe-refresh-seconds 7200 --heartbeat-seconds 900 --min-free-gb 10
Restart=on-failure
RestartSec=15
TimeoutStopSec=180
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
LimitNOFILE=65536
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-tick-compressor

[Install]
WantedBy=multi-user.target