# scripts/polymarket-bot.service
[Unit]
Description=Polymarket Mean-Reversion Market Maker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/home/botuser/polymarket-bot
Environment=PATH=/home/botuser/polymarket-bot/.venv/bin:/usr/bin:/bin
ExecStartPre=/bin/bash -c 'age -d -o /dev/shm/secrets/.env /home/botuser/polymarket-bot/.env.age || true'
ExecStart=/home/botuser/polymarket-bot/.venv/bin/python -m src.cli run --paper
ExecStopPost=/bin/bash -c 'rm -f /dev/shm/secrets/.env'

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/botuser/polymarket-bot/logs /dev/shm/secrets
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-bot

[Install]
WantedBy=multi-user.target
